# config/supervisor.conf
[program:openresty]
command=/opt/openresty/nginx/sbin/nginx -p /var/www/brainslug/ -c config/nginx.conf
autorestart=true
